# cookbooks/php/templates/default/php-fpm.init.erb
#!/bin/bash

# chkconfig: 345 60 30
# description: PHP FPM
# Source function library.
. /etc/init.d/functions

RETVAL=0

function phpstart {
        echo -n $"Starting PHP FPM"
        /usr/local/php/bin/php-cgi --fpm
        RETVAL=$?
        [ $RETVAL -eq 0 ] && success
        echo
        return $RETVAL

}

function phpstop {
        echo -n $"Stopping PHP FPM"
        kill -SIGINT `cat /usr/local/php/logs/php-fpm.pid`
        RETVAL=$?
        [ $RETVAL -eq 0 ] && success
        echo
        return $RETVAL

}

function phprestart {
	echo -n $"Restarting PHP FPM"
        kill -SIGUSR2 `cat /usr/local/php/logs/php-fpm.pid`
        RETVAL=$?
        [ $RETVAL -eq 0 ] && success
        echo
        return $RETVAL
}

case $1 in
start)
        cd /usr/local/php
        phpstart
        ;;
stop)
        cd /usr/local/php
        phpstop
        ;;
restart)
        cd /usr/local/php
	phprestart
        ;;
*)

        echo $"Usage: $0 {start|stop|restart}"
        exit 1
        ;;

esac

